# util/yahoo_scorebot.cron
#!/bin/bash

export DJANGO_SETTINGS_MODULE="project.prod_settings"
source /home/rhansen/.bashrc

${WORKON_HOME}/thorsboard/bin/python ${PROJECT_HOME}/thorsboard/acquire/yahoo/scores/scorebot.py -w current -u